# build/BRN_PI_Thermometer/DEBIAN/postinst
#bin/bash
sudo npm install --global http-server
chmod 644 /lib/systemd/system/BRNThermometer.service
systemctl daemon-reload
systemctl enable BRNThermometer.service
